# meta-photon/recipes-photon/photon-dashboard/files/photon-dashboard.service
[Unit]
Description=Photon Dashboard
After=systemd-udev-trigger.service
Conflicts=getty@tty1.service
DefaultDependencies=no
# Widen the start-limit window so a short retry interval doesn't blow
# through systemd's 5-restarts-in-10s gate and permanently fail the unit.
StartLimitBurst=20
StartLimitIntervalSec=30

[Service]
Type=simple
User=root
Environment=DISPLAY=:0 VK_ICD_FILENAMES=/usr/share/vulkan/icd.d/broadcom_icd.aarch64.json
WorkingDirectory=/root
ExecStart=/usr/bin/xinit /root/.xinitrc -- :0 vt1 -nolisten tcp -nocursor -s 0 -dpms
Restart=on-failure
RestartSec=500ms

[Install]
WantedBy=sysinit.target
